# scripts/com.messi.daily-git-pull.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.messi.daily-git-pull</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/offlee01/antigravity/find ssalmouk/ssalmeok-finding/scripts/daily_git_pull.sh</string>
    </array>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>9</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    <key>RunAtLoad</key>
    <false/>
    <key>StandardOutPath</key>
    <string>/tmp/openclaw/daily-git-pull-stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/openclaw/daily-git-pull-stderr.log</string>
</dict>
</plist>
